# templates/default/svscan.service.erb
[Service]
ExecStart=/usr/bin/svscanboot
Restart=always
